# deploy/birdlife-litestream.service
[Unit]
# Continuous offsite backup: replicate the shared bird DB to S3/B2 so an SD-card
# death doesn't lose the cottage's detection history. Needs the DB to exist, so
# it's ordered after the web app (whose db:prepare creates it). Only enabled when
# the LITESTREAM_* keys are set in .env (provision.sh skips it otherwise).
Description=Litestream replication of the bird DB to S3/B2
After=birdlife-web.service
Requires=birdlife-web.service

[Service]
User=pi
EnvironmentFile=/home/pi/birdlife/.env
ExecStart=/usr/bin/litestream replicate -config /home/pi/birdlife/deploy/litestream.yml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
